# data/autostart/typio.desktop.in
[Desktop Entry]
Version=1.0
Type=Application
Name=Typio
Comment=Start the Typio input method daemon
Exec=@TYPIO_HOST_DIR@/typio
TryExec=typio
Icon=typio-keyboard-symbolic
Terminal=false
StartupNotify=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
Categories=System;Utility;
Keywords=input;ime;keyboard;typio;
